# install/atm_bringup/share/atm_bringup/systemd/atm-slave-sentinel.service
[Unit]
Description=Sentinel externo de la Raspberry esclava del atomizador
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/atm/atm_slave_sentinel.py
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
